# data/redis_exporter.service
[Unit]
Description=Redis Exporter

[Service]
User=prometheus
ExecStart=/usr/bin/prometheus-redis-exporter

[Install]
WantedBy=default.target
